# tieba_signin.timer
[Unit]
Description=Baidu Tieba Automatic Signin Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target